# deploy/momentum-publish.service
[Unit]
Description=Publish bbterminal performance snapshot to the GitHub Pages repo
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=trader
Group=trader
WorkingDirectory=/opt/momentum-stack
Environment=MOMENTUM_DATA_DIR=/var/lib/momentum/data
# trader.env also sets: MOMENTUM_SITE_REPO_DIR, MOMENTUM_SITE_DATA_PATH,
# MOMENTUM_DEPLOY_KEY, MOMENTUM_PUBLISH_PUSH=1
EnvironmentFile=/etc/momentum/trader.env
ExecStart=/opt/momentum-stack/.venv/bin/momentum-publish --push

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/momentum
