# deploy/momentum-rebalance.timer
[Unit]
Description=Monthly momentum rebalance trigger

[Timer]
# Align this with bbterminal's next_rebalance_at (shown on the Performance page).
# Default: 06th of each month, 03:30 local — adjust to your strategy's schedule.
# The strategy-freshness gate protects you if the day drifts.
OnCalendar=*-*-06 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
