# deploy/momentum-trade-worker.service
[Unit]
Description=Momentum trade worker (holds IBKR + bbterminal credentials; the only thing that trades)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=trader
Group=trader
WorkingDirectory=/opt/momentum-stack
Environment=MOMENTUM_DB_PATH=/var/lib/momentum/momentum.db
Environment=MOMENTUM_DATA_DIR=/var/lib/momentum/data
# trader.env holds the secrets: IBIND_OAUTH1A_*, BBTERMINAL_*, SUPABASE_*, IBKR_ACCOUNT_ID,
# DRY_RUN (set false to allow live), REQUIRE_STRATEGY_HEALTHY, etc.
EnvironmentFile=/etc/momentum/trader.env
ExecStart=/opt/momentum-stack/.venv/bin/momentum-trade-worker --poll 5
Restart=on-failure
RestartSec=10

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/momentum

[Install]
WantedBy=multi-user.target
