# deploy/momentum-web.service
[Unit]
Description=Momentum web UI (read-only; holds no broker credentials)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=momentum
Group=momentum
WorkingDirectory=/opt/momentum-stack
Environment=MOMENTUM_WEB_HOST=127.0.0.1
Environment=MOMENTUM_WEB_PORT=8800
Environment=MOMENTUM_DB_PATH=/var/lib/momentum/momentum.db
Environment=MOMENTUM_DATA_DIR=/var/lib/momentum/data
# web.env must NOT contain IBKR/bbterminal credentials — this is the boundary.
EnvironmentFile=-/etc/momentum/web.env
ExecStart=/opt/momentum-stack/.venv/bin/momentum-web
Restart=on-failure
RestartSec=5

# Hardening: the web can only read data + read/write the SQLite job queue.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/momentum

[Install]
WantedBy=multi-user.target
