# deploy/momentum-mapping.timer
[Unit]
Description=Refresh the LEONTEQ+ETF -> IBKR mapping snapshot once a day

[Timer]
# The universe and IBKR listings change slowly; a daily resolve is plenty.
# 05:30 on weekdays, before the portfolio/rebalance jobs.
OnCalendar=Mon..Fri *-*-* 05:30:00
Persistent=true

[Install]
WantedBy=timers.target
