# deploy/momentum-strategies.timer
[Unit]
Description=Refresh the scheduled-strategies snapshot from bbterminal a few times a day

[Timer]
# Holdings only change on a rebalance, but the next-rebalance date and perf
# stats tick daily; a few refreshes a day keeps the page current and cheap.
OnCalendar=Mon..Fri *-*-* 06,12,18:00:00
Persistent=true

[Install]
WantedBy=timers.target
