# deploy/momentum-mapping.service
[Unit]
Description=Resolve LEONTEQ universe + ETFs to IBKR contracts for the read-only web /mapping page
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=trader
Group=trader
WorkingDirectory=/opt/momentum-stack
Environment=MOMENTUM_DATA_DIR=/var/lib/momentum/data
# trader.env holds the secrets: BBTERMINAL_*/SUPABASE_* (bbterminal admin API)
# and IBIND_OAUTH1A_*/IBKR_ACCOUNT_ID (IBKR OAuth), same file the worker uses.
EnvironmentFile=/etc/momentum/trader.env
ExecStart=/opt/momentum-stack/.venv/bin/momentum-mapping-snapshot

# Resolving ~1500 names one-by-one against IBKR takes a few minutes.
TimeoutStartSec=1800
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/momentum
