# deploy/momentum-portfolio.timer
[Unit]
Description=Refresh the IBKR portfolio snapshot every 15 min during market hours

[Timer]
# Weekday daytime (Europe), every 15 minutes. Adjust to your market/timezone.
OnCalendar=Mon..Fri *-*-* 08..22:00/15:00
Persistent=true

[Install]
WantedBy=timers.target
